# crontab
SHELL=/bin/bash

# AI.hab Production Crontab — supercronic format (standard cron syntax)
# All times UTC. PYTHONPATH e TZ herdados do container via supercronic.

# Hourly cycle (minuto 5 de cada hora)
5 * * * * /app/scripts/hourly_cycle.sh

# Daily update (07:00 UTC)
0 7 * * * /app/scripts/daily_update.sh

# Daily observation log (01:00 UTC — paper-trading observation period)
0 1 * * * cd /app && python scripts/daily_observation.py >> /app/logs/daily_observation.log 2>&1
