# docker/crontab
SHELL=/bin/bash

# AI.hab Production Crontab — supercronic format (standard cron syntax)
# All times UTC. PYTHONPATH e TZ herdados do container via supercronic.

# Hourly cycle — scoring completo (gates + regime + stops) — minuto 5 de cada hora
5 * * * * /app/scripts/hourly_cycle.sh >> /app/logs/hourly.log 2>&1

# Checagem rápida de stops — a cada 15min (minutos 20, 35, 50; não coincide com minuto 5)
20,35,50 * * * * cd /app && python3 scripts/check_stops.py >> /app/logs/stops.log 2>&1

# Daily update (07:00 UTC)
0 7 * * * /app/scripts/daily_update.sh

# Capital Manager snapshot — persiste capital_manager.json a cada 5 minutos
*/5 * * * * cd /app && python3 scripts/capital_snapshot.py --save >> /app/logs/capital_manager.log 2>&1

# ETH Bot 3 — ingestão spot + ciclo + stops (offset do BTC para não conflitar)
2 * * * * cd /app && python3 -m src.data.eth_ingestion >> /app/logs/eth_ingestion.log 2>&1
10 * * * * cd /app && python3 scripts/eth_hourly_cycle.py >> /app/logs/eth_bot3.log 2>&1
25,40,55 * * * * cd /app && python3 scripts/eth_check_stops.py >> /app/logs/eth_stops.log 2>&1
